# docker/supervisor/worker.conf
[program:worker]
command=php /usr/src/app/bin/console app:worker
autorestart=true
autostart=true