# docker/supervisor/conf.d/rr.conf
[program:rr]
command=/usr/src/app/./rr serve
autorestart=true
autostart=true
startsecs=5
stdout_logfile=/var/log/supervisor/output.log
stderr_logfile=/var/log/supervisor/output.log