# docker/supervisor/conf.d/worker.conf
[program:worker]
command=php /usr/src/app/bin/console app:worker
autorestart=true
autostart=true
stdout_logfile=/var/log/supervisor/output.log
stderr_logfile=/var/log/supervisor/output.log